# roles/wintergarten/files/wintergarten.service
[Unit]
Description=Wintergarten
After=syslog.target

[Service]
ExecStart=/usr/bin/wintergarten
Type=simple

[Install]
WantedBy=multi-user.target
