# templates/apps/supervisor/rhokthehood-beta.ini
[program:rhokthehood-beta]
command = /venvs/rhokthehood.com/beta/venv/bin/python /venvs/rhokthehood.com/beta/venv/bin/gunicorn --bind 127.0.0.1:8002 app:app
directory = /apps/rhokthehood.com/beta/
user = deploy
