# templates/apps/supervisor/rhokthehood.ini
[program:rhokthehood]
command = /venvs/rhokthehood.com/www/venv/bin/python /venvs/rhokthehood.com/www/venv/bin/gunicorn --bind 127.0.0.1:8002 app:app
directory = /apps/rhokthehood.com/www/
user = deploy
