# templates/apps/supervisor/wintergarten.ini
[program:wintergarten]
command = /venvs/citruspi.io/wintergarten-api/venv/bin/python /venvs/citruspi.io/wintergarten-api/venv/bin/gunicorn --bind 127.0.0.1:8000 server:app
directory = /apps/citruspi.io/wintergarten-api/
user = deploy
